# ansible/roles/install_controller_service/templates/controller.service.j2
[Unit]
Description={{ service_description }}
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ user }}
WorkingDirectory={{ home }}
ExecStart=docker compose -f {{ home }}/docker-compose-controller.base.yaml -f {{ home }}/docker-compose-controller.override.yaml up --wait
ExecStop=docker compose -f {{ home }}/docker-compose-controller.base.yaml -f {{ home }}/docker-compose-controller.override.yaml down

[Install]
WantedBy=multi-user.target
